# mkosi.extra/usr/share/factory/etc/profile.d/zfetch.sh
alias neofetch='/usr/bin/zfetch'
alias neowofetch='/usr/bin/zfetch'
alias fastfetch='/usr/bin/zfetch'
