# mkosi.profiles/liveiso-bootc-ostree/mkosi.extra/usr/lib/systemd/system/var-tmp.mount
[Unit]
Description=Larger tmpfs for /var/tmp on live system

[Mount]
What=tmpfs
Where=/var/tmp
Type=tmpfs
Options=size=50%%,nr_inodes=1m,x-systemd.graceful-option=usrquota

[Install]
WantedBy=local-fs.target
